# xandeum-pod-monitor.service
[Unit]
Description=Xandeum Pod Manager (Pod-Man) Web Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/xandeum-pod-monitor
ExecStart=/usr/bin/node /root/xandeum-pod-monitor/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xandeum-pod-monitor

# Security hardening
NoNewPrivileges=false
PrivateTmp=true

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
